# hack/backup/puck-etcd-backup.timer
[Unit]
Description=Run etcd snapshot nightly

[Timer]
OnCalendar=02:15
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
